# buildroot-ihost/package/debugger/debugger.service
[Unit]
Description=Coolkit debugger script

[Service]
Type=simple
ExecStart=node /home/root/debugger.js

[Install]
WantedBy=multi-user.target
